# websh.service
[Unit]
Description=websh mobile shell terminal
After=network.target

[Service]
# Must run as root: reads each user's ~/.config/websh.yaml and spawns shells as
# that user (drops privileges per session). Do NOT set NoNewPrivileges=true
# (breaks setuid) or ProtectHome (must read other users' home directories).
User=root
Type=simple
# The web UI is embedded in the binary; no static directory to deploy.
ExecStart=/usr/local/bin/websh \
  --bind 127.0.0.1:9631 \
  --secrets /var/lib/websh/secrets.json \
  --push-store /var/lib/websh/push_subs.json \
  --pam-service websh \
  --secure-cookie
Restart=on-failure
RestartSec=2
StateDirectory=websh
# /run/websh holds the web-session state file; preserve it across restarts so a
# restart doesn't log everyone out (cleared on reboot, which is fine).
RuntimeDirectory=websh
RuntimeDirectoryPreserve=yes

[Install]
WantedBy=multi-user.target
